# provisioning/ansible/roles/celery/templates/redis.ini.j2
[program:redis]
command=/usr/bin/redis-server
autorestart=true
stdout_logfile={{ ws_home }}/var/log/%(program_name)s.log
redirect_stderr=true
priority=10
